# dot_config/systemd/user/chezmoi-update.timer
[Unit]
Description=Hourly chezmoi update (git pull + apply)

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=10min
Persistent=true
Unit=chezmoi-update.service

[Install]
WantedBy=timers.target
